# cloud-custodian/tools/omnissm/init/systemd/omnissm-inventory.service
[Unit]
Description=Service for OmniSSM Inventory

[Service]
Type=simple
ExecStart=omnissm-agent ps
StandardError=journal
