# cloud-custodian/tools/omnissm/init/upstart/omnissm-inventory.cron
*/30 * * * * root /sbin/initctl start omnissm-inventory
