# cloud-custodian/tools/omnissm/init/upstart/omnissm-register.cron
0 6 * * * root /sbin/initctl start omnissm-register
@reboot root  /sbin/initctl start omnissm-register
