# cloud-custodian/tools/omnissm/init/upstart/omnissm-inventory.conf
description     "omnissm-inventory"
start on (runlevel [345] and started network)
stop on (runlevel [!345] or stopping network)

task

exec omnissm-agent ps > /var/log/omnissm-inventory.log 2>&1
